# templates/ffmpeg-monitor.service.template
[Unit]
Description=FFmpeg Monitor Web
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/monitor-web.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
